# vdsd.service.in
[Unit]
Description=vDS userspace daemon
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
ExecStart=@VDS_SYSTEMD_VDSD@
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=multi-user.target
